# systemd/welcome-framebuffer.service
[Unit]
Description=Display welcome pages on framebuffer
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/bin/bash /opt/welcome/scripts/display-welcome.sh
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
